# examples/kbr/kadmind.sv.conf
[program:kadmind]
command=/bin/bash -c "kadmind"
