# examples/kbr/kinit.sv.conf
[program:kadmin]
command=/usr/local/bin/kinit.sh
