# examples/kbr/krb5kdc.sv.conf
[program:kdc]
command=/bin/bash -c "krb5kdc"
