# examples/kbr/service.sv.conf
[program:service]
command=/usr/local/bin/start_service.sh
